# scripts/systemd/sidpro-web.service.example
# SIDPRO Web — systemd unit template (staging/production)
# Copy to /etc/systemd/system/sidpro-web.service and adjust paths.
# Do NOT commit real secrets — use EnvironmentFile outside the repository.

[Unit]
Description=SIDPRO Next.js Web
After=network.target sidpro-api.service
Wants=sidpro-api.service

[Service]
Type=simple
User=sidpro
Group=sidpro
WorkingDirectory=/opt/sidpro
EnvironmentFile=/etc/sidpro/sidpro.env
Environment=PORT=3000
ExecStart=/usr/bin/pnpm --filter @sidpro/web start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sidpro-web

# Hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
